# rootfs/qboxhd/rootfs/etc/inittab
# This is run first except when booting in single-user mode
#
::sysinit:/etc/init.d/rcS

# /bin/sh invocations on selected ttys
# Respawn a shell in the serial console 
ttyAS0::respawn:-/bin/sh

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/shutdown


